# root/etc/services.d/lidarr/run
#!/usr/bin/with-contenv bash

# cd /app/lidarr || exit

UMASK_SET=${UMASK_SET:-022}
umask "$UMASK_SET"

exec \
	s6-setuidgid abc /app/lidarr/Lidarr \
	-nobrowser -data=/config
