# root/etc/services.d/lidarr/run
#!/usr/bin/with-contenv bash

exec \
    s6-setuidgid abc /app/lidarr/bin/Lidarr \
    -nobrowser -data=/config
